# packages/buildroot/ouistiti/ouistiti.service
[Unit]
Description=An embedded web server for system management
After=syslog.target network.target

[Service]
Type=forking
PIDFile=/var/run/ouistiti.pid
ExecStartPre=/usr/bin/mkdir -p /var/run/ouistiti
ExecStart=/usr/sbin/ouistiti -D -p /var/run/ouistiti.pid -f /etc/ouistiti/ouistiti.conf
PrivateDevices=yes

[Install]
WantedBy=multi-user.target
